# infrastructure/systemd/fuel-maintenance.timer
# Nightly Maintenance Timer
# Triggers the nightly backup and maintenance routine at 3:00 AM

[Unit]
Description=NEXUS Nightly Maintenance Timer

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
