# infrastructure/systemd/sbc-client.service
# SBC Client Service
# Runs the SBC client software on each dispenser module

[Unit]
Description=NEXUS SBC Client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/nexus/sbc-client/sbc-client
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
